# .config/systemd/user/dunst.service
[Unit]
Description=Notification daemon

[Service]
ExecStart=/usr/bin/dunst
Type=simple

[Install]
WantedBy=default.target
