# .config/systemd/user/picom.service
[Unit]
Description=Compositor

[Service]
ExecStart=/usr/bin/picom
Type=simple

[Install]
WantedBy=default.target
