# .config/systemd/user/pulseaudio.service
[Unit]
Description=Audio server

[Service]
ExecStart=/usr/bin/pulseaudio
Type=simple

[Install]
WantedBy=default.target
