# .config/systemd/user/sxhkd.service
[Unit]
Description=Simple X hotkey daemon

[Service]
ExecStart=/usr/bin/sxhkd
Type=simple

[Install]
WantedBy=default.target
