# .config/systemd/user/weatheralertd.service
[Unit]
Description=Weather alert daemon

[Service]
ExecStart=/usr/local/bin/weatheralertd
Type=simple

[Install]
WantedBy=default.target
